# dot_config/systemd/user/webcam.timer
[Unit]
Description="run webcam service"

[Timer]
OnBootSec=30sec
OnUnitActiveSec=30sec

[Install]
WantedBy=timers.target
