# dot_config/systemd/user/webcam.service
[Unit]
Description="Configure webcam"

[Service]
Type=simple
ExecStart=%h/bin/configure-webcam.sh