# _config/systemd/user/skk-sync.timer
[Unit]
Description=Sync SKK dictionary

[Timer]
OnBootSec=1min
OnUnitActiveSec=1hour

[Install]
WantedBy=timers.target
