# _config/systemd/user/skk-sync.service
[Unit]
Description="sync SKK dictionary"

[Service]
Type=simple
ExecStart=%h/bin/skk-sync-dict
