# dot_templates/systemd/base-unit.service
[Unit]
Description={{_name_}}
Wants=
After=

[Service]
type=simple
ExecStart=

[Install]
WantedBy=default.target

{{_cursor_}}
